# com.sjtu.canvassync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sjtu.canvassync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>PROJECT_DIR/sync.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PROJECT_DIR</string>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>PROJECT_DIR/data/sync.log</string>

    <key>StandardErrorPath</key>
    <string>PROJECT_DIR/data/sync_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
